# recipes-core-luci/firewall/files/firewall.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2010 OpenWrt.org

ROOT_NFS=$(cat /proc/cmdline|grep 'root=/dev/nfs')
test -n "$ROOT_NFS" && exit 0

START=45

FW_LIBDIR=/lib/firewall

fw() {
	. $FW_LIBDIR/core.sh
	fw_$1
}

start() {
	fw start
}

stop() {
	fw stop
}

restart() {
	fw restart
}

reload() {
	fw reload
}
