# recipes-core-luci/gw-boot/files/boot-stage3
#!/bin/sh /etc/rc.common

START=40
STOP=98

scan_hardwares() {
	/lib/network/detect-modem.py
}

start() {
	scan_hardwares
	ifup wwan
	return 0
}
