# recipes-core-luci/gw-boot/files/runonce.service
[Unit]
Description=Pulsar Gateway Customer Run Once Function
After=network-online.target

[Service]
ExecStart=/etc/runonce/runonce.sh

[Install]
WantedBy=multi-user.target


